# softpack-mcp.service
[Unit]
Description=SoftPack MCP Server
After=network.target

[Service]
Type=exec
ExecStart=/usr/bin/make prod
Restart=always
RestartSec=10
User=ubuntu
WorkingDirectory=/home/ubuntu/work/softpack-mcp
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin:/bin:/usr/local/bin
Environment=PYTHONPATH=/home/ubuntu/work/softpack-mcp
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
